# _etc/systemd/system/media-games-nvme.automount
[Unit]
Description=Automount NVMe Games Drive

[Automount]
Where=/media/games/nvme
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target